# data/services/cpupower-gui-helper.service.in
[Unit]
Description=cpupower-gui system helper

[Service]
Type=dbus
BusName=org.rnd2.cpupower_gui.helper
ExecStart=@helperdir@/cpupower-gui-helper
