# data/services/cpupower-gui-user.service.in
[Unit]
Description=Apply cpupower-gui config at user login
After=default.target

[Service]
Type=oneshot
ExecStart=@bindir@/cpupower-gui config

[Install]
WantedBy=default.target
