# data/services/cpupower-gui.service.in
[Unit]
Description=Apply cpupower-gui config at boot
Documentation=https://github.com/vagnum08/cpupower-gui man:cpupower-gui(1)

[Service]
Type=oneshot
ExecStart=@bindir@/cpupower-gui config

[Install]
WantedBy=multi-user.target
