# data/services/cpupower-gui-helper.service.in
[Unit]
Description=cpupower-gui system helper
Documentation=https://github.com/vagnum08/cpupower-gui

[Service]
Type=dbus
BusName=org.rnd2.cpupower_gui.helper
ExecStart=@helperdir@/cpupower-gui-helper
